# perps-dashboard.service
[Unit]
Description=PerpsTrader Dashboard
After=network.target

[Service]
Type=simple
User=d
WorkingDirectory=/home/d/PerpsTrader
ExecStart=/home/d/.nvm/versions/node/v22.22.0/bin/node /home/d/PerpsTrader/bin/dashboard.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target