# perps-autoresearch.service
[Unit]
Description=PerpsTrader AutoResearch Bridge
After=network.target perps-redis.service
Wants=perps-redis.service

[Service]
Type=simple
User=d
WorkingDirectory=/home/d/autoresearch
ExecStart=/home/d/PerpsTrader/bin/research-engine/autoresearch-entry.js
Restart=on-failure
RestartSec=30
Environment=NODE_ENV=production
Environment=AUTORESEARCH_INTERVAL_MINUTES=60
Environment=AUTORESEARCH_AUTO_ADOPT_THRESHOLD=1.5
Environment=AUTORESEARCH_MAX_CONCURRENT=2
Environment=AUTORESEARCH_TIMEOUT_SECONDS=3600

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=perps-autoresearch

[Install]
WantedBy=multi-user.target
