# systemd/perps-news-global.service
[Unit]
Description=PerpsTrader Global Newsfeed Agent
After=network.target

[Service]
Type=simple
User=d
Group=d
WorkingDirectory=/home/d/PerpsTrader
ExecStart=/home/d/.nvm/versions/node/v22.22.0/bin/node /home/d/PerpsTrader/bin/news-agent.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=SEARCH_API_URL=http://localhost:8000/api/v1
Environment=PATH=/usr/bin:/home/d/.local/bin

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=perps-news

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
