# systemd/perps-dashboard.service
[Unit]
Description=PerpsTrader Dashboard
After=network.target perps-redis.service
Wants=perps-redis.service

[Service]
Type=simple
User=d
WorkingDirectory=/home/d/PerpsTrader
ExecStart=/usr/bin/node /home/d/PerpsTrader/bin/dashboard/dashboard-server.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target